# board/sbnb/sbnb/rootfs-overlay/usr/lib/systemd/system/sbnb-tunnel.service
[Unit]
Description=sbnb tunnel service
After=sbnb.service
Requires=sbnb.service
ConditionPathExists=/mnt/sbnb/tunnel-start.sh

[Service]
Type=simple
ExecStart=/bin/sh /mnt/sbnb/tunnel-start.sh
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target

