# board/sbnb/sbnb/rootfs-overlay/usr/lib/systemd/system/sbnb-watchdog.timer
[Unit]
Description=sbnb watchdog timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
