# board/sbnb/sbnb/rootfs-overlay/usr/lib/systemd/system/sbnb-watchdog.service
[Unit]
Description=sbnb watchdog

[Service]
Type=oneshot
ExecStart=/usr/bin/sbnb-watchdog.sh
TimeoutStartSec=50s
